# config/com.highfocus.bot.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.highfocus.bot.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/iliaburkhovetskii/highfocus/venv/bin/python</string>
        <string>/Users/iliaburkhovetskii/highfocus/bot.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/iliaburkhovetskii/highfocus</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>NetworkState</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/iliaburkhovetskii/highfocus/bot.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/iliaburkhovetskii/highfocus/bot_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/iliaburkhovetskii/highfocus/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>UserName</key>
    <string>iliaburkhovetskii</string>
    
    <key>GroupName</key>
    <string>staff</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>

